# backend/stonks.service
Description=Stonks App Server
After=syslog.target
After=network.target

[Service]
Type=simple
User=lev
WorkingDirectory=/home/lev/app/stock_price_prediction_app/backend
ExecStart=/home/lev/app/stock_price_prediction_app/backend/server
RestartSec=10
Restart=always

[Install]
WantedBy=multi-user.target